# crontab
# Write out status every five minutes to ensure curator is working as expected
*/5 * * * * echo "curator running on $(hostname)" >> /dev/stdout
